# deploy/signaling.service
[Unit]
Description=WebRTC PortMap Signaling Server (+embedded TURN)
Documentation=https://github.com/your-org/webrtc-portmap
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=webrtc-portmap
Group=webrtc-portmap
WorkingDirectory=/var/lib/webrtc-portmap

# 所有启动参数集中放在 env 文件里（systemd 会按空白拆分 $SIGNALING_ARGS 为多个参数）
EnvironmentFile=/etc/webrtc-portmap/signaling.env
ExecStart=/usr/local/bin/signaling $SIGNALING_ARGS

Restart=on-failure
RestartSec=3s
# 优雅退出：signaling 收到 SIGTERM 会关闭 TURN/HTTP
KillSignal=SIGTERM
TimeoutStopSec=10s

# ---- 适度安全加固 ----
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# 仅允许写数据目录（signaling.json 落盘在这里）
ReadWritePaths=/var/lib/webrtc-portmap
# 若要监听 80/443 等特权端口，取消下一行注释（默认 8443/3478 无需特权）
# AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
